# build/one/supervisord-one.conf
# Supervisord configuration for Amgix-One (all-in-one container)

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# Redirect supervisord logs to stdout/stderr
[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# RabbitMQ - Message queue (priority 1 - starts first)
[program:rabbitmq]
command=/app/log-wrapper.sh "rabbit" /usr/sbin/rabbitmq-server
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

# Qdrant - Vector database (priority 2 - starts second)
[program:qdrant]
command=/app/log-wrapper.sh "qdrant" /usr/local/bin/qdrant --config-path /app/config.yaml
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

# Amgix API (priority 13 - waits for dependencies, runs on port 8234)
[program:api]
command=/app/log-wrapper.sh "amgix-api" /app/wait-for-deps.sh /app/run-api.sh
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=10
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

# Amgix Encoder 1 (priority 11 - waits for dependencies)
[program:encoder1]
command=/app/log-wrapper.sh "amgix-encoder1" /app/wait-for-deps.sh /app/run-encoder.sh
environment=HOSTNAME=encoder1,AMGIX_MODEL_FAMILY="all"
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=11
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

# Amgix Encoder 2 (priority 12 - waits for dependencies)
[program:encoder2]
command=/app/log-wrapper.sh "amgix-encoder2" /app/wait-for-deps.sh /app/run-encoder.sh
environment=HOSTNAME=encoder2,AMGIX_MODEL_FAMILY="all"
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=12
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
